# etc/zsh/profile.d/01-env-vars.zsh
set -a
source "${ENVY_ROOT}/etc/env"
set +a