# etc/zsh/profile.d/90-prompt.zsh
# shellcheck shell=zsh

export STARSHIP_CONFIG="${ENVY_ROOT}/etc/starship/starship.toml"
eval "$(starship init zsh)"