# etc/bash/profile.d/00-shell.sh
# shellcheck shell=bash

case "$-" in
  *i*)
    set -o vi
    ;;
esac
