# etc/bash/profile.d/63-go.sh
# shellcheck shell=bash

export GOPATH="${HOME}/.go"
